# extensions/68kmlaorg/macproxy_plus.service
# https://github.com/hunterirving/macproxy_plus

[Unit]
Description=MacProxy Plus, an extensible HTTP proxy that connects early computers to the Internet.
After=network.target

[Service]
WorkingDirectory=/opt/macproxy_plus
# Point PATH at the venv’s bin *first*, then fall back to the system
Environment=PATH=/opt/macproxy_plus/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin
# Launch the proxy with the venv’s python
ExecStart=/opt/macproxy_plus/venv/bin/python /opt/macproxy_plus/proxy.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
